# crontab
# METARMap Crontab Configuration
# This crontab runs the METARMap every 5 minutes from 7 AM to 7 PM,
# and turns off the lights at 8 PM.
# Assumes the project is installed in /home/pi/metarmap
# For custom schedules, visit https://crontab.guru/

# Run METARMap every 5 minutes from 7:00 AM to 7:00 PM
*/5 7-18 * * * cd /home/pi/metarmap && source metarmap_env/bin/activate && sudo python3 metar.py

# Turn off lights at 8:00 PM
0 20 * * * cd /home/pi/metarmap && source metarmap_env/bin/activate && sudo python3 pixelsoff.py
